# meta-mitysom-5csx-vdk/recipes/init-ifupdown/init-ifupdown-1.0/if-up.d/start_udhcpd
#!/bin/sh
echo HELLO
echo $IFACE $MODE
if [ $IFACE = usb0 ] && [ $MODE = start ]
then
	echo "starting udhcpd"
	touch /tmp/udhcpd.leases
	/usr/sbin/udhcpd /etc/udhcpd.conf
fi
